# infra/lxc/mantis-control.service
[Unit]
Description=Mantis-DNS control plane
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=mantis
Group=mantis
EnvironmentFile=/etc/mantis-control/mantis-control.env
WorkingDirectory=/opt/mantis-dns/app
ExecStart=/opt/mantis-dns/venv/bin/uvicorn mantis_control.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/mantis-dns
# Holds MANTIS_SIGNING_KEY_PATH (signing_key.bin) — deliberately outside
# /opt/mantis-dns/app, which install-rocky.sh wipes and recreates from a
# fresh checkout on every reinstall. A signing key stored there would
# regenerate on every reinstall, silently invalidating every already-running
# filter node's cached public key until it's manually restarted.
ReadWritePaths=/etc/mantis-control

[Install]
WantedBy=multi-user.target
